# contrib/daemon/gentoo-init/init.d/lmsd
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header$

extra_started_commands="reload"

start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --pidfile ${LMSD_PIDFILE} --exec ${LMSD_EXEC} -- -h ${LMSD_DBHOST} \
	    -d ${LMSD_DBNAME} -u ${LMSD_DBUSER} -p ${LMSD_DBPASS} -H ${LMSD_HOSTNAME} -P ${LMSD_PIDFILE} > /dev/null 2>&1
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop  --pidfile ${LMSD_PIDFILE}
	eend $?
}

reload() {
	einfo "Reloading all modules ${SVCNAME} ..."
	${LMSD_EXEC} -h ${LMSD_DBHOST} -d ${LMSD_DBNAME} -u ${LMSD_DBUSER} -p ${LMSD_DBPASS} -H ${LMSD_HOSTNAME} -q
	einfo "Reloaded all modules ${SVCNAME} ..." 
	eend $?
}
